# supervisord.d/shibd.ini
[program:shibd]
command=/usr/sbin/shibd -F

