# supervisord.d/sssd.ini
[program:sssd]
command=/usr/sbin/sssd -i
